# scripts/pai-sync/pai-sync-daemon.service
[Unit]
Description=PAI Memory Sync Daemon
Documentation=file:///home/ec2-user/repos/imladris/docs-site/docs/specs/memory-sync.md
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/repos/imladris/scripts/pai-sync

# Run via bun directly (compile to binary for production later via bun build)
ExecStart=/home/ec2-user/.bun/bin/bun run daemon.ts

Restart=always
RestartSec=5
WatchdogSec=60

# Postgres connection URL — set by Ansible from BWS, never hardcoded
EnvironmentFile=/etc/pai-sync/env

# Inherit bun path
Environment=PATH=/home/ec2-user/.bun/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/ec2-user

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pai-sync

[Install]
WantedBy=multi-user.target
